# fileconverter.conf
[program:FileConverterBot]
directory=/home/ubuntu/FileConverterBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/FileConverterBot/logfile_err.log
stdout_logfile=/home/ubuntu/FileConverterBot/logfile.log